# systemd/vpn-dns-sync.path
[Unit]
Description=Watch for VPN IP changes on tun0

[Path]
PathModified=/proc/net/fib_trie

[Install]
WantedBy=multi-user.target
