# systemd/vpn-client-ip-dns-sync@.path
[Unit]
Description=Watch for VPN IP changes on %i

[Path]
PathModified=/proc/net/fib_trie

[Install]
WantedBy=multi-user.target
