# systemd/vpn-dns-sync.service
[Unit]
Description=Sync VPN IP to DigitalOcean DNS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-dns-sync.sh
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
